# daemons/s3sync.service
[Unit]
Description=Sync AWS S3 with Local

[Service]
ExecStart=/usr/bin/s3sync.sh

[Install]
WantedBy=multi-user.target
